# shell/profile.d/path.sh
# Begin /etc/profile.d/path.sh
export PATH=/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin:/bin:/sbin:$PATH
export PKG_CONFIG_PATH=/usr/lib64/pkgconfig:/usr/lib/pkgconfig:/usr/local/lib/pkgconfig
# End /etc/profile.d/path.sh
